# Developpement-BeagleBoneBlack/WithUSBDongle/drone-agent.service
[Unit]
Description= Launch drone agent pairingprocess

[Service]
Type=simple
ExecStart=/usr/bin/drone-agent.sh

[Install]
WantedBy=multi-user.target